# bin/deploy_vnc
#!/bin/sh

# READY

######################################################################
# Script version is yymmdd-HHMMSS in UTC, date -u +%y%m%d-%H%M%S
######################################################################
SCRIPTVERSION=150618-090257

pkg install -y x11/xorg x11-wm/xfce4-wm \
Thunar xfce4-tumbler \
www/epiphany \
net/tightvnc
#x11-wm/openbox
#www/midori

if [ ! -e /usr/local/bin/vncserver ]; then
    echo "ERROR: install failed"
    exit 1
fi

mkdir ~/.vnc

if [ ! -f ~/.vnc/xstartup ]; then
    cat >~/.vnc/xstartup <<EOF
#!/bin/sh
xrdb \$HOME/.Xresources &
xfwm4 &
xfce4-panel &
xterm -bg black -fg gray &
EOF
    chmod +x ~/.vnc/xstartup
fi

cat >/usr/local/etc/rc.d/vnc <<EOF
#!/bin/sh
# Copyright (c) 2015, John Ko
# PROVIDE: vnc
# REQUIRE: FILESYSTEMS LOGIN NETWORKING
. /etc/rc.subr

name="vnc"
start_cmd="my_start"
stop_cmd="my_stop"
my_start()
{
    /usr/local/bin/vncserver -depth 24 -geometry 1280x800 -nolisten tcp :1 &
}
my_stop()
{
    # TODO need a more elegant solution than loop to kill
    local LOOP=1
    while [ \$LOOP -eq 1 ]; do
        if  /bin/pgrep -lf vncserver > /dev/null 2> /dev/null ; then
            /bin/pkill -lf vncserver > /dev/null 2> /dev/null
        else
            local LOOP=0
        fi
        sleep 1
    done
}
load_rc_config \$name
run_rc_command "\$1"
EOF
chmod +x /usr/local/etc/rc.d/vnc

cat >/tmp/vncserver.patch.b64 <<EOF
LS0tIHZuY3NlcnZlcgkyMDEzLTEwLTA5IDE2OjEwOjI1LjAwMDAwMDAwMCArMDAw
MAorKysgdm5jc2VydmVyLm5ldwkyMDEzLTEwLTA5IDE2OjEyOjQ1LjAwMDAwMDAw
MCArMDAwMApAQCAtMjk2LDcgKzI5Niw3IEBAIHN1YiBDaGVja0Rpc3BsYXlOdW1i
ZXIKIAogICAgIHNvY2tldChTLCAkQUZfSU5FVCwgJFNPQ0tfU1RSRUFNLCAwKSB8
fCBkaWUgIiRwcm9nOiBzb2NrZXQgZmFpbGVkOiAkIVxuIjsKICAgICBldmFsICdz
ZXRzb2Nrb3B0KFMsICZTT0xfU09DS0VULCAmU09fUkVVU0VBRERSLCBwYWNrKCJs
IiwgMSkpJzsKLSAgICB1bmxlc3MgKGJpbmQoUywgcGFjaygnUyBuIHgxMicsICRB
Rl9JTkVULCA2MDAwICsgJG4pKSkgeworICAgIHVubGVzcyAoYmluZChTLCBzb2Nr
YWRkcl9pbig2MDAwICsgJG4sICZJTkFERFJfQU5ZKSkpIHsKIAljbG9zZShTKTsK
IAlyZXR1cm4gMDsKICAgICB9CkBAIC0zMDQsNyArMzA0LDcgQEAgc3ViIENoZWNr
RGlzcGxheU51bWJlcgogCiAgICAgc29ja2V0KFMsICRBRl9JTkVULCAkU09DS19T
VFJFQU0sIDApIHx8IGRpZSAiJHByb2c6IHNvY2tldCBmYWlsZWQ6ICQhXG4iOwog
ICAgIGV2YWwgJ3NldHNvY2tvcHQoUywgJlNPTF9TT0NLRVQsICZTT19SRVVTRUFE
RFIsIHBhY2soImwiLCAxKSknOwotICAgIHVubGVzcyAoYmluZChTLCBwYWNrKCdT
IG4geDEyJywgJEFGX0lORVQsIDU5MDAgKyAkbikpKSB7CisgICAgdW5sZXNzIChi
aW5kKFMsIHNvY2thZGRyX2luKDU5MDAgKyAkbiwgJklOQUREUl9BTlkpKSkgewog
CWNsb3NlKFMpOwogCXJldHVybiAwOwogICAgIH0K
EOF

cat /tmp/vncserver.patch.b64 | openssl enc -d -base64 >/tmp/vncserver.patch

cd /usr/local/bin
patch </tmp/vncserver.patch

sysrc dbus_enable="YES"

service dbus start

cat >~/.vnc/passwd.b64 <<EOF
29g8/XJ6FFg=
EOF
chmod 400 ~/.vnc/passwd.b64
cat ~/.vnc/passwd.b64 | openssl enc -d -base64 >~/.vnc/passwd
chmod 400 ~/.vnc/passwd

vncserver -depth 24 -geometry 1280x800 -nolisten tcp :1

echo
JAILIP=`ifconfig | grep inet | awk '{print $2}'`
echo "Connect with ssh -L 5901:${JAILIP}:5901 you@your-server-ip"
echo "Then connect with password to vnc://localhost:5901"
